# supervisord.conf
[supervisord]
nodaemon=true

[program:api-watcher]
autorestart=true
command=/bin/sh -c 'npm run api-watcher'

[program:docs-server]
autorestart=true
command=/bin/sh -c 'npm run docs-server'
